# deploy/gft-sync-shifts.timer
[Unit]
Description=Sync GFT shifts hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
