# deploy/gft-sync-screenings.service
[Unit]
Description=GFT Usher - Sync Screenings

[Service]
Type=oneshot
User=gft
Group=gft
WorkingDirectory=/opt/gft-usher-tools
Environment=PATH=/home/gft/.bun/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/gft/.bun/bin/bun run src/server/sync-cli.ts screenings
